# service/rpmbuild/BUILD/swap-monitor/swap-monitor.service
[Unit]
Description=Swap monitoring service

[Service]
Type=simple
ExecStart=/usr/bin/swap-monitor.sh
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target
